# systemd/apt-upgrade/apt-upgrade.service
[Unit]
Description=Run full system upgrade

[Service]
Type=oneshot
ExecStart=/usr/local/bin/apt-upgrade-service