# systemd/template/template.timer
[Unit]
Description=Run service at scheduled time

[Timer]
OnCalendar=*-*-* 07:00

[Install]
WantedBy=timers.target