# systemd/template/template.service
[Unit]
Description=

[Service]
Type=oneshot
ExecStart=/usr/local/bin/template-service