# systemd/node-reboot/node-reboot.service
[Unit]
Description=Just do a node reboot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/reboot-service