# systemd/node-reboot/node-reboot.timer
[Unit]
Description=Run service every morning

[Timer]
OnCalendar=*-*-* 07:00

[Install]
WantedBy=timers.target