# systemd/node-poweroff/node-poweroff.timer
[Unit]
Description=Run service every night

[Timer]
OnCalendar=*-*-* 03:10


[Install]
WantedBy=timers.target