# systemd/node-poweroff/node-poweroff.service
[Unit]
Description=Set wakeup timer and poweroff

[Service]
Type=oneshot
ExecStart=/usr/local/bin/node-poweroff-service 300